# nutlogger.cron
# Example cron job for nutlogger with flock
# Launches nutlogger if not running, checks every five minutes
*/5 * * * * flock -n /home/nutlogger/nutlogger.lock /home/nutlogger/nutlogger/nutlogger
